# groups/main/research-papers/systemd/nanoclaw-qa-heal.service
[Unit]
Description=NanoClaw Paper Q&A auto-healer
Documentation=https://github.com/JungwooHur/nanoclaw
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/home/jw/nanoclaw/.env
WorkingDirectory=/home/jw/nanoclaw
ExecStart=/usr/bin/python3 /home/jw/nanoclaw/groups/main/research-papers/auto_fix_qa.py
ExecStart=/usr/bin/python3 /home/jw/nanoclaw/groups/main/research-papers/auto_save_qa.py
StandardOutput=journal
StandardError=journal
