# groups/main/research-papers/systemd/paperclaw-qa-heal.service
[Unit]
Description=PaperClaw Paper Q&A auto-healer
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/paperclaw/.env
WorkingDirectory=%h/paperclaw
# The `-` prefix makes each step's failure non-fatal so a transient error in
# one (e.g. a Notion HTTP 429) doesn't abort the rest — without it a failed
# auto_fix_qa.py blocks the downstream healers (a real past incident).
ExecStart=-/usr/bin/python3 %h/paperclaw/groups/main/research-papers/auto_fix_qa.py
ExecStart=-/usr/bin/python3 %h/paperclaw/groups/main/research-papers/auto_save_qa.py
ExecStart=-/usr/bin/python3 %h/paperclaw/groups/main/research-papers/collect_papers.py --dedupe
# Structurally enforce the paper-body cleanup the agent doesn't reliably do:
# strip translated back-matter (References/Acknowledgements) + leaked source URLs
# (NotebookLM image URLs, ar5iv citation URLs). Recently-edited pages only, so the
# 5-minute cadence stays cheap.
ExecStart=-/usr/bin/python3 %h/paperclaw/groups/main/research-papers/heal_paper_pages.py
StandardOutput=journal
StandardError=journal
