# systemd/paperclaw-watchdog.service
[Unit]
Description=PaperClaw WhatsApp connection watchdog
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/paperclaw/scripts/paperclaw-watchdog.sh
StandardOutput=journal
StandardError=journal
